# src/etc/services.d/pmmp/run
#!/usr/bin/with-contenv sh

cd /minecraft

exec s6-setuidgid minecraft ./start.sh --no-wizard



