# etc/services.d/pmmp/run
#!/usr/bin/execlineb -P

with-contenv
s6-envuidgid minecraft

cd /minecraft
/usr/local/bin/php ./PocketMine-MP.phar --no-wizard
